# systemd/system/tcsd.service.d/trousers-changer.conf
[Service]
ExecStartPre=/sbin/trousers_changer_migrate
ExecStartPre=/sbin/trousers_changer_identify
